# scripts/rc.local
#!/bin/bash

dhcpcd
sleep 1s
./home/root/usb_ether.sh

exit 0
